# Phase2/zm/etc/zonemgr.service
[Unit]
Description=ZoneManager
After=network.target

[Service]
ExecStart=/usr/lib/zonemgr/zonemanager.py --database sqlite:////var/lib/zonemgr/zones.db server -e network 127.0.0.1:9009

[Install]
WantedBy=multi-user.target

