# Phase 1/dns-zone-mgr/etc/crontab
5 0 * * * /usr/lib/zonemgr/zonemanager.py --database sqlite:////var/lib/zonemgr/zones.db resign
